# backup.timer
[Unit]
Description=Backup SD card

[Timer]
OnCalendar=daily
 
[Install]
WantedBy=timers.target
